# ble-mqtt-bridge-ui.service
[Unit]
Description=BLE to MQTT Bridge Web UI
After=network.target ble-mqtt-bridge.service

[Service]
Type=simple
User=ble-mqtt
Group=ble-mqtt
WorkingDirectory=/opt/ble-mqtt
ExecStart=/opt/ble-mqtt/venv/bin/python /opt/ble-mqtt/bridge_ui.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ble-mqtt
PrivateTmp=true

[Install]
WantedBy=multi-user.target
